# fc-tinder.service
[Unit]
Description=FC Tinder Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/fc_tinder
Environment=NODE_ENV=production
Environment=PORT=80
Environment=CSV_DIR=good csvs
ExecStart=/usr/bin/node /home/ubuntu/fc_tinder/dist/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
